# vision.conf
start on runlevel [2345]
stop on runlevel [016]
respawn

script
  sleep 10
  exec /home/ubuntu/2383/2016-Vision/2016-Vision
end script
